# config/systemd/loki.service
[Unit]
Description=Loki log aggregation
After=network.target

[Service]
ExecStart=/usr/local/bin/loki -config.file=/etc/loki/loki-config.yaml
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
